# carini/examples/msr/c3c2/c3_works_24hrs/crontab
* * * * * /home/idardik/projects/compositional_ii/tla-sym-assump-synth/examples/msr/naive_decomp/monitor_proc.sh
